# trance-daemon/debian/postinst
#!/bin/sh
# postinst script for trance-daemon
set -e

if [ "$1" = "configure" ]; then
    echo "=========================================================="
    echo "  Setting up Trance Screensaver Daemon..."
    echo "=========================================================="

    # Enable the service globally for all user sessions using standard deb-systemd-helper
    if command -v deb-systemd-helper >/dev/null; then
        deb-systemd-helper --user unmask trance-daemon.service || true
        deb-systemd-helper --user enable trance-daemon.service || true
    fi

    # Start/restart the service for all currently logged-in users using systemd's user manager API
    if command -v loginctl >/dev/null && command -v systemctl >/dev/null; then
        loginctl list-users --no-legend | while read -r uid user linger state; do
            if [ -n "$user" ]; then
                echo "-> Restarting/Starting service for: $user"
                systemctl --user --machine="$user@" daemon-reload || true
                systemctl --user --machine="$user@" restart trance-daemon.service || true
            fi
        done
    fi
    echo "  Setup complete. Trance Daemon is now running!"
    echo "=========================================================="
fi

exit 0
